# RHCSA/ Shell Scripting Basics/crontab
0 17 * * * /bin/bash /home/yassen/mysql_backup.sh
